# files/etc/crontabs/wifisongadmin
*/2 * * * *	/usr/lib/wifixi/watchdog
0   4 * * * reboot
